# cron/crontab
*/15 * * * * /opt/stats.sh && cp /opt/stats.json /var/log/nginx/
